# misc/install/nats/nats-service.service
[Unit]
Description=NATS service

[Service]
User=root
ExecStart=nats-server -m 8222
Restart=always

[Install]
WantedBy=multi-user.target

